# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-postgresql/finish
#!/bin/sh
echo "Stopping PostgreSQL gracefully..."
su postgres -c "pg_ctl stop -D /data/pg_data -m fast" || true
